# myquiz.service
[Unit]
Description=MyQuiz Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/html
ExecStart=/usr/bin/node server/index.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=5001
Environment=JWT_SECRET=your-super-secret-jwt-key-change-this-in-production

[Install]
WantedBy=multi-user.target 